# clickerd/etc/clicker-cycle.service
[Unit]
Description=Clicker Cycle

[Service]
Type=simple
ExecStart=/usr/local/bin/clicker cycle -l 10 -r

[Install]
WantedBy=multi-user.target
