# clickerd/etc/clicker-http.service
[Unit]
Description=Clicker HTTP

[Service]
Type=simple
ExecStart=/usr/local/bin/clicker-http

[Install]
WantedBy=multi-user.target
